# launchd/com.nanoclaw.refresh-claude-creds.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nanoclaw.refresh-claude-creds</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/milagpt/nanoclaw/nanoclaw-v2/scripts/refresh-claude-creds.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/milagpt/nanoclaw/nanoclaw-v2/logs/refresh-claude-creds.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/milagpt/nanoclaw/nanoclaw-v2/logs/refresh-claude-creds.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
